# linuxmini/systemd/user/dunst-config-watch.path
[Unit]
Description=Watch Dunst config for changes

[Path]
PathModified=/home/rash/.config/dunst/dunstrc

[Install]
WantedBy=graphical-session.target
