# linuxmini/systemd/user/dunst.service
[Unit]
Description=Start Dunst notification daemon
After=graphical-session.target
PartOf=graphical-session.target
# Add condition to ensure Hyprland is running
ConditionPathExists=%t/hypr

[Service]
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/dunst -conf /home/rash/.config/dunst/dunstrc
# Only show notification if marker file exists, then remove it
ExecStartPost=/bin/bash -c 'if [ -f "${XDG_RUNTIME_DIR}/dunst/restarted" ]; then sleep 0.5 && dunstify "Dunst has been restarted 🎉" --icon=dialog-information --timeout=3000 && rm "${XDG_RUNTIME_DIR}/dunst/restarted"; fi'
Restart=on-failure
RestartSec=1

[Install]
WantedBy=graphical-session.target
