# linuxmini/systemd/user/hypr-keybinds-watcher.path
[Unit]
Description=Watch for Hyprland keybinds_source.conf changes
Documentation=man:systemd.path(5)

[Path]
PathModified=%h/.config/hypr/keybinds_source.conf
Unit=hypr-keybinds-translate.service

[Install]
WantedBy=default.target