# linuxcommon/systemd/user/sops-secrets.service
[Unit]
Description=Decrypt sops secrets and inject into user environment
Documentation=https://github.com/getsops/sops
After=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%h/.dotfiles/common/scripts/decrypt-secrets.sh

[Install]
WantedBy=default.target
