# systemd/user/dunst-reload.path
[Unit]
Description=Watch dunst.config for changes

[Path]
PathChanged=/home/rash/.config/dunst/dunst.config

[Install]
WantedBy=default.target
