# systemd/user/kanata-status-listener.service
[Unit]
Description=Kanata Status Listener
Documentation=https://github.com/jtroo/kanata
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=/home/rash/.local/bin/kanata-tools listen
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
Environment=HOME=/home/rash
WorkingDirectory=/home/rash/.config/scripts/kanata

# Security restrictions
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp
PrivateTmp=no

# Allow access to /tmp for status files and kanata config for persistent state
ReadWritePaths=/tmp /home/rash/.config/kanata

[Install]
WantedBy=default.target